# deploy/ticker-calendar.service
[Unit]
Description=Ticker Calendar Alert Scheduler
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# EDIT: set your Linux username and project path before enabling
User=roop
Group=roop
WorkingDirectory=/home/roop/Desktop/ticker-calender
Environment=NTFY_TOPIC=test_myalerts
Environment=TICKER_CALENDAR_LOG_DIR=/home/roop/Desktop/ticker-calender/logs
ExecStart=/home/roop/Desktop/ticker-calender/.venv/bin/python run_server.py serve
Restart=always
RestartSec=10
# Give yfinance + ntfy time to finish if a job is running during shutdown
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
