# deploy/ticker-calendar.service
[Unit]
Description=Ticker Calendar Alert Scheduler
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# EDIT: set your Linux username and project path before enabling
User=roop
Group=roop
WorkingDirectory=/home/roop/Desktop/ticker-calender
# REQUIRED: set a private, random ntfy topic (e.g. `openssl rand -hex 16`).
# ntfy.sh topics are public — a guessable name leaks alerts and allows spoofing.
Environment=NTFY_TOPIC=CHANGE_ME_TO_A_RANDOM_PRIVATE_TOPIC
Environment=TICKER_CALENDAR_LOG_DIR=/home/roop/Desktop/ticker-calender/logs
ExecStart=/home/roop/Desktop/ticker-calender/.venv/bin/python run_server.py serve
Restart=always
RestartSec=10
# Give yfinance + ntfy time to finish if a job is running during shutdown
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
